# tools/swat-memory/hooks/com.c3.swat-memory.maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.c3.swat-memory.maintenance</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>-m</string>
    <string>swat_memory.maintenance</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>__REPO__/src</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>15</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__HOME__/.claude/swat-memory/maintenance.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.claude/swat-memory/maintenance.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
